# mail/templates/certbot/certbot.service.j2
[Install] 
WantedBy=multi-user.target

[Unit]
Description=Certbot renewal service

[Service]
ExecStart=/usr/bin/docker run \
    --attach STDOUT \
	--attach STDERR \
	--network {{ mail_bridge_name }} \
	--rm \
	--cap-drop all \
	--volume {{ mail_vol_base }}/certbot:/etc/letsencrypt \
	--volume {{ mail_vol_base }}/ovh.ini:/etc/letsencrypt/ovh.ini:ro \
	certbot/dns-ovh \
	renew

ExecStartPost=/usr/bin/docker restart dovecot
ExecStartPost=/usr/bin/docker restart postfix
