# nextcloud/templates/nextcloud-backup.service.j2
[Install] 
WantedBy= multi-user.target

[Unit]
Description=Nextcloud backup service

[Service]
ExecStart={{ nextcloud_backup_script }}