# data/system/systemd/fail2ban.service.d/override.conf
[Service]
# filesystem access
PrivateDevices=yes
PrivateTmp=yes
ProtectHome=true
ProtectSystem=strict
ProtectProc=invisible
ReadWritePaths=-/var/spool/postfix/maildrop

# service directories
RuntimeDirectory=fail2ban
LogsDirectory=fail2ban
StateDirectory=fail2ban
ConfigurationDirectory=fail2ban

# kernel
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true

# privileges
CapabilityBoundingSet=CAP_AUDIT_READ CAP_DAC_READ_SEARCH CAP_NET_ADMIN CAP_NET_RAW

# misc
RestrictNamespaces=true
SystemCallArchitectures=native
RestrictRealtime=true
MemoryDenyWriteExecute=true
LockPersonality=true
ProtectHostname=true
RemoveIPC=true
RestrictSUIDSGID=true
ProtectClock=true